# fridge-hfdc-curl.service
[Unit]
Description=Fridge tracker uploader

[Service]
ExecStart=/home/pi/fridge-hfdc/fridgeUpload.sh

[Install]
WantedBy=multi-user.target
